# dist-gimp/root/.config/autostart/gimp.desktop
[Desktop Entry]
Encoding=UTF-8
Version=0.9.4
Type=Application
Name=GIMP
Comment=
Exec=/usr/bin/gimp
StartupNotify=false
Terminal=false
Hidden=false

